# install/lora.cron
* *    * * *    root    /opt/lorawan-sk-gateway/bin/lora_to_cloud > /dev/null 2>&1
